# ulyaoth-haproxy/SOURCES/haproxy.service
[Unit]
Description=HAProxy - The Reliable, High Performance TCP/HTTP Load Balancer
Documentation=https://www.haproxy.org/
After=network.target remote-fs.target nss-lookup.target

[Service]
Type=forking
PIDFile=/run/haproxy.pid
ExecStart=/usr/local/sbin/haproxy -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s QUIT $MAINPID
PrivateTmp=true
LimitNOFILE=200000

[Install]
WantedBy=multi-user.target